# docker/sphinxsearch/supervisor/conf.d-avaliable/sphinx-updateIndexes.conf
[program:sphinx-updateRtIndex]
command=bash /usr/local/bin/updateIndexes.sh
environment=
autostart=false
autorestart=false
numprocs=1
redirect_stderr=true
priority=1100
